# scripts/systemd/bulennode-server.service
[Unit]
Description=BulenCoin Node (server-full)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bulen
WorkingDirectory=/opt/bulencoin/bulennode
EnvironmentFile=-/etc/default/bulennode-server
ExecStart=/usr/bin/env npm start
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target

